# overlay/etc/inittab
::sysinit:/sbin/setup.sh
ttyS0::respawn:-/bin/sh
::ctrlaltdel:/sbin/reboot
